# runtime/openclaw/systemd/openclaw-telegram-bot.service
[Unit]
Description=OpenClaw Telegram Bot local-first
After=network-online.target openclaw-gateway.service
Wants=network-online.target
ConditionPathExists=/srv/tesis/repo/runtime/openclaw/wrappers/openclaw-telegram-bot.sh

[Service]
Type=simple
EnvironmentFile=/etc/tesis-os/openclaw.env
ExecStart=/bin/bash -lc '/srv/tesis/repo/runtime/openclaw/wrappers/openclaw-telegram-bot.sh'
Restart=always
RestartSec=5
User=openclaw
Group=openclaw
WorkingDirectory=/srv/tesis/repo
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/srv/tesis/repo /etc/tesis-os/openclaw.env /etc/tesis-os/domains
ReadWritePaths=/var/lib/herramientas/openclaw /var/cache/herramientas/openclaw /var/log/openclaw /srv/tesis/intercambio/openclaw
StandardOutput=append:/var/log/openclaw/openclaw-telegram-bot.log
StandardError=append:/var/log/openclaw/openclaw-telegram-bot.log

[Install]
WantedBy=multi-user.target
