# config/systemd/tesis-healthcheck.service
[Unit]
Description=Healthcheck del Sistema Operativo de Tesis
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/tesis/repo
ExecStart=/bin/bash /srv/tesis/repo/ops/actualizacion/healthcheck_tesis.sh
